# services/otr/run
#!/bin/bash

echo "Starting OTR Script"
if [ -z ${PUID+x} ]
then
  exec ${WORKDIR}/otr-auto.sh
else
  exec /sbin/setuser abc ${WORKDIR}/otr-auto.sh
fi
